# debian/postinst
#!/bin/sh -e

depmod

systemctl enable  xilinx-apf-driver.service
systemctl start   xilinx-apf-driver.service

